# php7.4-fpm/supervisord.conf
[supervisord]
logfile = /dev/null
loglevel = info
nodaemon=true

[program:ssh]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:cron]
command = cron -f -L 15
autostart=true
autorestart=true

[program:php-fpm]
command = php-fpm
autostart = true
autorestart = true
stderr_logfile=/var/log/php-fpm/error.log
stdout_logfile=/var/log/php-fpm/out.log
